# server/supervisor.conf
[program:rtt-server]
command=xrootd -c /etc/xrootd/xrootd.cfg -l /var/log/xrootd/xrootd.log
user=xrootd
autorestart=true
